# docker/supervisord.dev.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

# Angular dev server (ng serve with watch)
[program:angular]
command=/bin/sh -c "cd /app/frontend && npm start"
environment=PATH="/app/frontend/node_modules/.bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
startsecs=10
stopasgroup=true
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/angular.log
stderr_logfile=/var/log/supervisor/angular_error.log

# Spring Boot (mvn spring-boot:run with watch)
[program:backend]
command=bash -c "cd /app/backend && ./mvnw spring-boot:run -Dspring-boot.run.arguments='--spring.profiles.active=h2'"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/backend.log
stderr_logfile=/var/log/supervisor/backend_error.log
environment=SPRING_PROFILES_ACTIVE=h2,JAVA_OPTS="-Xmx512m -Xms256m -agentlib:jdwp=transport=dt_socket,server=y,suspend=n,address=*:5005"

# Nginx reverse proxy
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx_error.log

[unix_http_server]
file=/var/run/supervisor.sock

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
